# crontab/get-crontab-file.sh
#!/bin/bash
sudo mkdir origin
sudo cp /etc/crontab origin/
