# crontab/get-crontab-file.sh
#!/bin/bash
sudo cp /etc/crontab origin/
